# dom0/qvm-remote-dom0.service
[Unit]
Description=qvm-remote-dom0 — VM-to-dom0 command bridge
After=qubesd.service

[Service]
Type=simple
ExecStart=/usr/bin/qvm-remote-dom0
EnvironmentFile=-/etc/qubes/remote.conf
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
# Not enabled by default.  Use: qvm-remote-dom0 enable
WantedBy=multi-user.target
